# versions/6.1.0/share/init/dv.conf.erb
description "This is an upstart job file for DataVirtualization"

pre-start script
bash << "EOF"
  mkdir -p /var/log/dv
  chown -R dv /var/log/dv
EOF
end script

start on started network-services
stop on stopped network-services
respawn

limit nofile 4096 4096

script
bash << "EOF"
  su - dv
  /opt/dv/jboss/bin/standalone.sh <%= server_opts %> >> /var/log/dv/dv.log 2>&1
EOF
end script
